# files/system/etc/profile.d/settheming.sh
gsettings set org.gnome.desktop.interface gtk-theme 'Tokyonight-Dark-Compact' 
gsettings set org.gnome.desktop.interface color-scheme 'prefer-dark'
gsettings set org.gnome.desktop.interface icon-theme 'Flat-Remix-Blue-Dark'
gsettings set org.gnome.desktop.interface cursor-theme 'Bibata-Modern-Ice'
gsettings set org.gnome.desktop.interface cursor-size 20